# answer-server-final-package/answer.service
[Unit]
Description=Answer Service with All Customizations
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/answer_server
ExecStart=/root/answer_server/answer-linux run -C /root/answer_server/data
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
